# _site_setup/upstart_pwc.conf
start on started mysql
stop on stopped mysql

respawn
respawn limit 3 5

setuid paul
setgid paul

exec /home/paul/pwc/_site_setup/service_demon_pwc.sh