# _site_setup/upstart_pwc.conf
start on started mysql
stop on stopped mysql

respawn
respawn limit 3 5

setuid pwc
setgid pwc

exec /home/paul/pwc/_site_setup/service_demon_pwc.sh